# helper/com.funnel.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.funnel.helper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/funnel-helper</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>10240</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>10240</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/var/log/funnel-helper.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/funnel-helper.log</string>
</dict>
</plist>
